# deploy/medoraai-gateway-9000.service
# MedoraAI Monitoring Gateway — port 9000, TCP/HL7 -> Django API
# O'rnatish: sudo cp deploy/medoraai-gateway-9000.service /etc/systemd/system/ && sudo systemctl daemon-reload && sudo systemctl enable --now medoraai-gateway-9000

[Unit]
Description=MedoraAI Monitoring Gateway (port 9000)
After=network.target medoraai-backend-8001.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/medoraai
Environment=PATH=/root/medoraai/backend/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/root/medoraai
ExecStart=/root/medoraai/backend/venv/bin/python -m uvicorn monitoring_gateway.main:app --host 0.0.0.0 --port 9000 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
